# supervisor/conf.d/tmule.conf
[program:tmule]
command=/usr/local/bin/tmule server
process_name=%(program_name)s 
directory=/home/turtlebot
user=turtlebot
autostart=true
autorestart=true 
redirect_stderr=true

